# packaging/debian/camelmailer-worker.service
[Unit]
Description=CamelMailer delivery worker (queue, DKIM, tracking, webhooks)
Documentation=https://github.com/camelmailer/camelmailer
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
User=camelmailer
Group=camelmailer
EnvironmentFile=-/etc/camelmailer/env
Environment=CAMELMAILER_CONFIG_FILE_PATH=/etc/camelmailer/camelmailer.yml
ExecStart=/usr/bin/camelmailer worker
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
